# lbindex-dist
#!/bin/bash
# chkconfig: 2345 99 01
# description: Esquema de inicialização do LBIndex
# processname: lbindex

. /usr/local/lbneo/virtenvlb2.6/src/LBIndex/lbindexmg
